# supervisord.conf
[program:hubot]
command = /home/dotcloud/current/node_modules/.bin/coffee bin/hubot-wrapper
directory = /home/dotcloud/current
